# app.service
[Service]
ExecStart=/usr/bin/java -jar /home/ubuntu/app.jar

[Install]
WantedBy=multi-user.target
